# conf/systemd/pfwd.service
[Unit]
Description=pfwd Service

[Service]
Type=simple
User=root
Group=root
PIDFile=/run/pfwd/pfwd.pid
RuntimeDirectory=pfwd

TimeoutStartSec=60
TimeoutStopSec=60

ExecStart=/usr/sbin/pfwd
ExecStop=/bin/kill -s TERM $MAINPID

Restart=always
RestartSec=300s

StandardOutput=syslog
StandardError=syslog

LimitNOFILE=131072

[Install]
WantedBy=multi-user.target
Alias=pfwd.service
